# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:bridge-default]
command=bash -c "mkdir -p /data/store && exec /app/whatsapp-bridge"
directory=/app/bridge-default
environment=PORT="7481",WHATSAPP_BRIDGE_API_KEY="%(ENV_WHATSAPP_BRIDGE_API_KEY)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=10
startsecs=5

[program:bridge-second]
command=bash -c "mkdir -p /data/store2 && exec /app/whatsapp-bridge"
directory=/app/bridge-second
environment=PORT="7482",WHATSAPP_BRIDGE_API_KEY="%(ENV_WHATSAPP_BRIDGE_API_KEY)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=10
startsecs=5

[program:gateway]
command=python /app/main.py
directory=/app
environment=MCP_PORT="8080",MCP_AUTH_TOKEN="%(ENV_MCP_AUTH_TOKEN)s",MCP_BASE_URL="%(ENV_MCP_BASE_URL)s",GATEWAY_DB_PATH="/data/gateway.db",WHATSAPP_DB_PATH="/data/store/messages.db",WHATSAPP_BRIDGE_URL="http://localhost:7481",WHATSAPP_BRIDGE_API_KEY="%(ENV_WHATSAPP_BRIDGE_API_KEY)s",OAUTH_STATE_PATH="/data/oauth_state.json",RESULTS_DIR="/data/results"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=10
startsecs=3
